# deploy/systemd/user/dragonfly-feed-cache.service
[Unit]
Description=Dragonfly Flash feed cache refresher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_DIR}}
Environment=PYTHONUNBUFFERED=1
EnvironmentFile={{ENV_FILE}}
ExecStart=/usr/bin/python3 {{PROJECT_DIR}}/dragonfly_telegram_poster.py \
  --env-file {{ENV_FILE}} \
  --dragonfly-account backup_1 \
  --request-delay 2 \
  refresh-feed-cache-watch --count 60 --offset 0 --interval 20
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=dragonfly-bridge.target
